# raspberry-dashboard.service
[Unit]
Description=Raspberry Dashboard
After=graphical.target

[Service]
Type=simple
User=alarm
Environment=SNAPCAST_HOST=127.0.0.1:1705
Environment=DASHBOARD_CONFIG=/home/alarm/config.toml
ExecStart=/home/alarm/raspberry-dashboard
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
